# crontab
44 20 * * * python3 /Users/nedimdrekovic/Library/Mobile\ Documents/com\~apple\~CloudDocs/Python/automatic_tu_sport_script_one_browser.py Laufsport 1
